# autostart/kanshi.desktop
[Desktop Entry]
Exec=/usr/bin/kanshi
Name=kanshi
Type=Application
