# autostart/minidlnad.desktop
[Desktop Entry]
Exec=/usr/bin/minidlnad -f /home/batoshu/.config/minidlna/minidlna.conf
Name=minidlnad
Type=Application
