# autostart/transmission-daemon.desktop
[Desktop Entry]
Exec=/usr/bin/transmission-daemon
Name=transmission-daemon
Type=Application
