# autostart/udiskie.desktop
[Desktop Entry]
Exec=/usr/bin/udiskie
Name=udiskie
Type=Application
